# launchers/com.emma.icu-job-search-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.emma.icu-job-search-dashboard</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/aidanhutchison/.local/bin/emma-job-dashboard-start</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>
  <key>WorkingDirectory</key>
  <string>/Users/aidanhutchison/Documents/New project 3</string>
  <key>StandardOutPath</key>
  <string>/Users/aidanhutchison/Library/Logs/emma-icu-job-search-dashboard.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/aidanhutchison/Library/Logs/emma-icu-job-search-dashboard.err.log</string>
</dict>
</plist>
